# src/main/resources/scripts/crontab
SHELL=/bin/sh
MAILTO=admin@puyallupfamilyhistorycenter.org

0 3 * * * "cd $HOME/Dropbox/Install/FHC && ./launch-server.sh >>/tmp/fhc/logs/cronrun 2>&1"
